# config/etf.cron
6 8 * * * root /bin/find /opt/omd/sites/etf/var/nagios/archive/* -mtime +7 -exec rm {} \;
10 8 * * sun root /bin/find /var/cache/ncgx/inventory* -mtime +30 -exec rm {} \;

